# hardening/crontab
# m h  dom mon dow   command

0 * * * * python /home/njonji/Desktop/ASTROCYTECH/AstroCode/hardening/backup_script.py